# docker/postfix/config/supervisord/postfix.conf
[program:postfix]
directory=/etc/postfix
command=/usr/libexec/postfix/master -c /etc/postfix -d
priority=3
autostart=true
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
