# docker/postfix/config/supervisord/rsyslog.conf
[program:rsyslog]
command=/usr/sbin/rsyslogd -n
priority=1
autorestart=true
autostart=true
stdout_logfile=/var/log/rsyslog.log
stderr_logfile=/var/log/rsyslog.log
